# raspie-camera.service
[Unit]
Description=Raspberry Pi Camera and Audio Capture Service
After=network.target

[Service]
Type=simple
User=dawg
WorkingDirectory=/home/dawg/raspberry_pie_camera_capture
ExecStart=/home/dawg/raspberry_pie_camera_capture/start-daily-recording.sh
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Set environment variables
Environment="DISPLAY=:0"
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target 