# config/imx296-camera.service
[Unit]
Description=IMX296 Global Shutter Camera Capture Service
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/raspberry_pie_camera_capture
ExecStartPre=/home/pi/raspberry_pie_camera_capture/bin/restart_camera.sh
ExecStart=/home/pi/raspberry_pie_camera_capture/bin/run_imx296_capture.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target 