# setup/imx296-camera-monitor.service
[Unit]
Description=IMX296 Camera Service with Status Monitor
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/raspberry_pie_camera_capture
ExecStart=/usr/bin/python3 /home/pi/raspberry_pie_camera_capture/bin/start_camera_with_monitor.py --monitor --no-output
Environment=PYTHONPATH=/home/pi/raspberry_pie_camera_capture
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Service management
KillMode=mixed
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target 